# systemd/security-mailer.timer
[Unit]
Description=Run security-mailer once a day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
